# .config/autostart/bauh.desktop
[Desktop Entry]
Type=Application
Name=bauh (tray)
Exec=/path/to/bauh-tray

